# devops/roles/setup_service/files/spring-boot-app.service
[Unit]
Description=spring-boot-app

[Service]
ExecStart=/usr/local/bin/my_service

[Install]
WantedBy=multi-user.target