# devops/roles/setup_service/files/spring-boot-app.service
[Unit]
Description=spring-boot-app

[Service]
ExecStart=java -jar /opt/spring-boot-app.war

[Install]
WantedBy=multi-user.target